# buildroot/board/rockchip/common/base/etc/profile.d/pppd.sh
if [ -e /dev/ttyUSB2 ]; then
                source /etc/pppd.sh &
fi
